# roles/docker/templates/etc/systemd/system/docker.service.d/custom-docker-opts.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H unix:// {{ docker_service_systemd_custom_opts }}
