# roles/docker/templates/etc/systemd/system/docker-cleanup.service
[Unit]
Description=Docker System Cleanup
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune -a --volumes -f