# roles/docker/templates/etc/systemd/system/podman-cleanup.service
[Unit]
Description=Podman System Cleanup
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=/usr/bin/podman system prune -a --volumes -f